# molly-wizard.service
[Unit]
Description=Molly Gateway Wizard
After=network.target tailscaled.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/molly-gateway
ExecStart=/usr/bin/python3 /home/pi/molly-gateway/wizard.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install] 
WantedBy=multi-user.target
