# molly-wizard.service
[Unit]
Description=Molly Gateway Wizard
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/molly-gateway
ExecStart=/usr/bin/python3 /home/pi/molly-gateway/wizard.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
